# ci/ci/supervisord.conf
[program:webserver]
command={{ PYTHON }} {{ BASE_DIR }}/manage.py runserver 5656