# ci/ci/configs/supervisor/leo---test.com.conf
[program:leo---test.com.local]
directory = /home/zdimon/tmp/leo---test.com/pressa-besa/backend
command = /home/zdimon/Desktop/work/pressa-besa/backend/venv/bin/python3 manage.py runserver 0.0.0.0:8081
autostart = true
autorestart = true
environment=PYTHONPATH="/home/zdimon/Desktop/work/pressa-besa/backend/venv"
stderr_logfile=/home/zdimon/Desktop/work/pressa-besa/ci/logs/leo---test.com.local.log